# deploy/supervisor.ini
[program:enlighted_customers_cherrypy]
user=marc
command=/home/marc/envs/enlighted/bin/python /home/marc/python_test_apps/rma_customers/customers.py
directory=/home/marc/python_test_apps/rma_customers/
autostart=false
autorestart=true
logfile=log.txt
